# deploy/hl-factura.service
[Unit]
Description=HL Factura - Portal CFDI 4.0 multi-despacho
After=network.target mysql.service

[Service]
Type=simple
User=hlfactura
WorkingDirectory=/opt/hl-factura
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=COOKIE_SECURE=1
ExecStart=/usr/bin/npm start -- -p 3000
Restart=always
RestartSec=5
# Los certificados y XMLs viven en /opt/hl-factura/data
ReadWritePaths=/opt/hl-factura/data

[Install]
WantedBy=multi-user.target
